# deploy/systemd/weathermetars-web.service.template
[Unit]
Description=Weather METARs web server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
EnvironmentFile=__APP_ENV_FILE__
Environment=NODE_ENV=production
ExecStart=__BUN_BIN__ .output/server/index.mjs
Restart=always
RestartSec=5
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=__APP_DATA_DIR__

[Install]
WantedBy=multi-user.target
